# abraxas.service
[Unit]
Description=ABRAXAS dynamic color temperature daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/abraxas --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
